# router/devhost-router.service.tmpl
[Unit]
Description=Devhost Wildcard Router (FastAPI)
After=network.target

[Service]
Type=simple
# Replace {{USER}} with the target system user who will run the service
User={{USER}}
# Working directory where the router lives
WorkingDirectory=/home/{{USER}}/devhost/router
# Path to the uvicorn or python binary inside the project's venv. Replace {{VENV_BIN}}.
ExecStart={{VENV_BIN}} uvicorn app:app --host 127.0.0.1 --port 7777
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

# Usage instructions (example):
# 1) Fill placeholders: replace {{USER}} and {{VENV_BIN}} (e.g. /home/you/.venv/bin/uvicorn)
# 2) Copy to /etc/systemd/system/devhost-router.service
# 3) Enable and start: sudo systemctl enable --now devhost-router.service