# router/devhost-router.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>devhost.router</string>

  <key>ProgramArguments</key>
  <array>
    <!-- Point this to the venv uvicorn binary or to system uvicorn. Replace {{VENV_BIN}} accordingly. -->
    <string>{{VENV_BIN}}</string>
    <string>app:app</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>5555</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/{{USER}}/devhost/router</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/devhost-router.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/devhost-router.err</string>
</dict>
</plist>

<!-- Usage:
  1) Replace {{USER}} and {{VENV_BIN}} (e.g. /Users/you/.venv/bin/uvicorn)
  2) Copy to ~/Library/LaunchAgents/devhost.router.plist
  3) load: launchctl load ~/Library/LaunchAgents/devhost.router.plist
  4) unload: launchctl unload ~/Library/LaunchAgents/devhost.router.plist
-->